# scripts/gcs-daemon.service
[Unit]
Description=GCS Admin AI Daemon
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/gcs
ExecStart=/root/.nvm/versions/node/v20.18.1/bin/node /var/www/gcs/scripts/gcs-daemon.js
Restart=always
RestartSec=5
Environment=GCS_DAEMON_PORT=9876
EnvironmentFile=/var/www/gcs/.env.production

[Install]
WantedBy=multi-user.target
